# .config/autostart/xmodmap.desktop
[Desktop Entry]
Type=Application
Exec=sh -c "xmodmap ~/.xmodmaprc"
